# files/shared/systemd/user/chezmoi-sync.timer
[Unit]
Description=Chezmoi auto sync timer

[Timer]
OnBootSec=5m
OnUnitInactiveSec=1h

[Install]
WantedBy=timers.target
