# files/shared/systemd/user/chezmoi-sync.service
[Unit]
Description=Chezmoi auto sync dotfiles

[Service]
Type=oneshot
ExecStart=/usr/bin/chezmoi re-add
